# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info
user=root

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:postgres]
command=/usr/lib/postgresql/15/bin/postgres -D /data/postgres
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/postgres.log
stderr_logfile=/var/log/postgres.log
priority=1

[program:redis]
command=redis-server --appendonly yes --dir /data/redis
autostart=true
autorestart=true
stdout_logfile=/var/log/redis.log
stderr_logfile=/var/log/redis.log
priority=2

[program:rabbitmq]
command=/bin/bash -c 'RABBITMQ_NODENAME=rabbit@localhost rabbitmq-server'
autostart=true
autorestart=true
stdout_logfile=/var/log/rabbitmq.log
stderr_logfile=/var/log/rabbitmq.log
priority=3
environment=HOME="/var/lib/rabbitmq",RABBITMQ_NODENAME="rabbit@localhost",RABBITMQ_MNESIA_BASE="/var/lib/rabbitmq/mnesia",RABBITMQ_LOG_BASE="/var/log/rabbitmq",RABBITMQ_USE_LONGNAME="false"

[program:kamism]
command=/app/kamism-server
directory=/app
autostart=false
autorestart=true
stdout_logfile=/var/log/kamism.log
stderr_logfile=/var/log/kamism.log
priority=10
startsecs=3

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx-access.log
stderr_logfile=/var/log/nginx-error.log
priority=20

[program:init]
command=/entrypoint.sh
autostart=true
autorestart=false
stdout_logfile=/var/log/init.log
stderr_logfile=/var/log/init.log
priority=30
startsecs=0
